# supervisord/postgresql.conf
[program:postgresql]
command=/usr/bin/postgres -c config_file="/etc/postgresql/postgresql.conf"
user=postgres
